# doorbell.service

[Unit]
Description=doorbell

[Service]
User=
Group=
ExecStart= 

[Install]
WantedBy=sysinit.target
